# app/templates/projects/distribution/src/standard/scripts/api.service
[Unit]
Description=<%= slug %>

[Service]
WorkingDirectory=/opt/api/
ExecStart=/opt/api/scripts/api.sh

[Install]
WantedBy=multi-user.target
